# image/overlays/etc/systemd/system/autodarts-webpanel.service
[Unit]
Description=Autodarts Pi OS local web panel
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/autodarts-webpanel
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
