# image/overlays/etc/systemd/system/autodarts-kiosk.service
[Unit]
Description=Autodarts Pi OS kiosk display
After=network-online.target autodarts-firstboot.service
Wants=network-online.target

[Service]
Type=simple
User=autodarts
Group=autodarts
Environment=HOME=/home/autodarts
Environment=XDG_RUNTIME_DIR=/tmp/runtime-autodarts
ExecStartPre=+/usr/bin/touch /var/log/autodarts-kiosk.log
ExecStartPre=+/usr/bin/chown autodarts:autodarts /var/log/autodarts-kiosk.log
ExecStart=/usr/local/bin/autodarts-kiosk
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
