# image/overlays/etc/systemd/system/autodarts-kiosk.service
[Unit]
Description=Autodarts Pi OS kiosk display
After=systemd-user-sessions.service autodarts-firstboot.service autodarts-webpanel.service
Wants=autodarts-webpanel.service
Conflicts=getty@tty1.service

[Service]
Type=simple
User=autodarts
Group=autodarts
Environment=HOME=/home/autodarts
Environment=XDG_RUNTIME_DIR=/tmp/runtime-autodarts
Environment=XINITRC=/usr/local/bin/autodarts-kiosk-session
WorkingDirectory=/home/autodarts
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=no
StandardInput=tty
StandardOutput=journal
StandardError=journal+console
ExecStartPre=+/usr/bin/touch /var/log/autodarts-kiosk.log
ExecStartPre=+/usr/bin/chown autodarts:autodarts /var/log/autodarts-kiosk.log
ExecStartPre=+/usr/bin/install -d -m 0755 -o autodarts -g autodarts /home/autodarts
ExecStartPre=+/usr/bin/install -d -m 0700 -o autodarts -g autodarts /home/autodarts/.cache /home/autodarts/.config /home/autodarts/.local/state/autodarts-pi-os /tmp/runtime-autodarts
ExecStartPre=+/usr/bin/chown -R autodarts:autodarts /home/autodarts/.cache /home/autodarts/.config /home/autodarts/.local /tmp/runtime-autodarts
ExecStart=/usr/local/bin/autodarts-kiosk
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
