# image/overlays/etc/systemd/system/getty@tty1.service.d/autodarts-autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin autodarts --noclear --noissue %I $TERM
